# com.gexbot.bookmap.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gexbot.bookmap.update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/cooljl31/Documents/Bookmap cloud notes/run_gex_update.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>180</integer>

    <key>StandardOutPath</key>
    <string>/Users/cooljl31/Documents/Bookmap cloud notes/gex_launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/cooljl31/Documents/Bookmap cloud notes/gex_launchd_error.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
